# redis/config/cluster/systemd/cluster.service
[Unit]
Description=Redis Node %I
[Service]
Type=simple
ExecStart=redis-server /app/node_%I/redis.conf